# supervisor/supervisord.conf
[unix_http_server]
file=/run/ossm/supervisor.sock
chmod=0700

[supervisord]
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/ossm/supervisord.pid
nodaemon=true
minfds=1024
user=ato

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/ossm/supervisor.sock

[program:xvnc]
priority=10
command=/opt/ossm/scripts/start-vnc.sh
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:xfce]
priority=20
command=/opt/ossm/scripts/start-xfce.sh
autorestart=true
startsecs=3
environment=DISPLAY=":1",HOME="/home/ato",XDG_RUNTIME_DIR="/run/ossm/xdg"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:websockify]
priority=30
command=/opt/ossm/scripts/start-websockify.sh
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:gateway]
priority=40
command=/usr/bin/python3 /opt/ossm/gateway/lab_gateway.py
autorestart=true
startsecs=2
environment=OSSM_WORKSPACE="/foss/designs"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:nginx]
priority=50
command=/opt/ossm/scripts/start-nginx.sh
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

